# sway-start.service
[Unit]
Description=Configure Sway for eGPU

[Service]
Type=oneshot
ExecStart=/usr/bin/sway_setup.sh

[Install]
WantedBy=graphical.target
